# cms-clusters/wordpress/app/wordpress_single_apache_test/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock  

[supervisord]
logfile=/tmp/supervisord.log 
logfile_maxbytes=50MB        
logfile_backups=10           
loglevel=info                
pidfile=/var/run/supervisord.pid 
nodaemon=true            
minfds=1024                  
minprocs=200                

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock 

[program:nginx]
command=apache2ctl start
autostart=true
startsecs=10
autorestart=true
startretries=3

[program:php-fpm]
command=php-fpm7.0
autostart=true
startsecs=10
autorestart=true
startretries=3

[program:mysql]
command=mysql_startup.sh
autostart=true
startsecs=10
autorestart=true
startretries=3

[program:initwordpress]
command=install_wordpress.sh
autostart=true
startsecs=2
autorestart=unexpected
startretries=3

[include]
files = /etc/supervisor/config.d/*.ini

