# thingsboard/docker/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock  

[supervisord]
logfile=/tmp/supervisord.log 
logfile_maxbytes=50MB        
logfile_backups=10           
loglevel=info                
pidfile=/var/run/supervisord.pid 
nodaemon=true            
minfds=1024                  
minprocs=200                

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock 

[program:thingsboard]
command=/usr/share/thingsboard/bin/thingsboard.jar
autostart=true
startsecs=10
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile = /var/log/supervisor/thingsboard.log
logfile_maxbytes=50MB
logfile_backups=5

#[program:postgres]
#command=init_postgresql.sh

[include]
files = /etc/supervisor/config.d/*.ini
